# queue_manager_server.socket
[Unit]
Description=Queue Manager gunicorn socket

[Socket]
ListenStream=/run/queue_manager_gunicorn.sock

[Install]
WantedBy=sockets.target
